# v0.2/script/goAuto
#! /bin/sh
# /etc/init.d/goAuto

### BEGIN INIT INFO
# Provides:          goAuto
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple script to start a program at boot
# Description:       A simple script from www.stuffaboutcode.com which will start / stop a program a boot / shutdown.
### END INIT INFO

# If you want a command to always run, put it here

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting the TechLab Brain"
    # run application you want to start
    python /home/pi/Documents/tagNFC_DB_server/v0.2/script/theBrain.py >/tmp/brainLog_normal.txt 2>/tmp/brainLog_error.txt &
    echo "Starting telegram"
    python /home/pi/Documents/tagNFC_DB_server/v0.2/script/telegramBot.py >/tmp/botLog_normal.txt 2>/tmp/botLog_error.txt &
    ;;
  stop)
    echo "Stopping theBrain!"
    # kill application you want to stop
    kill -9 $(ps aux | grep '[p]ython /home/pi/Documents/tagNFC_DB_server/theBrain/theBrain.py' | awk '{print $2}')
    echo "Stopping telegram"
    kill -9 $(ps aux | grep '[p]ython /home/pi/Documents/tagNFC_DB_server/theBrain/theBot.py' | awk '{print $2}')
    ;;
  *)
    echo "Usage: /etc/init.d/theBrain {start|stop}"
    exit 1
    ;;
esac

exit 0
